# src/playbooks/roles/rtags/files/rtags-rdm.conf
description "rtags rdm service"
author  "km45"

start on runlevel [2345]
stop on runlevel [016]

script
    exec su vagrant --command 'rdm'
end script

respawn
